# packaging/systemd/shorewalld.service
[Unit]
Description=shorewall-nft monitoring + DNS-set API daemon
Documentation=file:///usr/share/doc/shorewall-nft/shorewalld.md
# Start after shorewall-nft so the ruleset is present when
# shorewalld probes for tables.
After=network.target shorewall-nft.service
Wants=network.target

[Service]
Type=simple
# Default config: scrape every netns under /run/netns plus own ns,
# Prometheus on :9748, dnstap endpoint OFF.
# Override on the command line or via
# /etc/systemd/system/shorewalld.service.d/override.conf.
ExecStart=/usr/bin/shorewalld --listen-prom :9748 --netns auto
Restart=on-failure
RestartSec=5s

# Security hardening (matches shorewall-nft.service where applicable).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/shorewalld /var/lib/shorewalld /var/log
PrivateTmp=yes

# Capabilities:
#   CAP_NET_ADMIN  — setns + nft list table from other namespaces
#   CAP_SYS_PTRACE — needed to enter other netns via setns(2)
#   CAP_NET_RAW    — conservative, not strictly required
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_PTRACE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_PTRACE

# Create the dnstap socket dir with the right mode up front. The
# daemon will unlink/rebind the socket on start.
RuntimeDirectory=shorewalld
RuntimeDirectoryMode=0750
StateDirectory=shorewalld
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
