# packaging/systemd/shorewalld@.service
[Unit]
Description=shorewall-nft monitoring daemon for netns %i
Documentation=file:///usr/share/doc/shorewall-nft/shorewalld.md
# The instance name %i maps to a netns under /run/netns.
# Use this template if you prefer one daemon per namespace instead
# of letting a single daemon serve all of them (shorewalld.service
# with --netns auto).
After=network.target shorewall-nft@%i.service
Wants=network.target

[Service]
Type=simple
# Per-instance Prometheus port: 9748 + netns-hash-offset would be
# cleaner, but we keep it simple with a fixed port; operators that
# want multi-instance scraping should use per-instance NetworkNamespacePath
# drop-ins to bind inside the target netns and pick distinct ports.
EnvironmentFile=-/etc/sysconfig/shorewalld
EnvironmentFile=-/etc/default/shorewalld
ExecStart=/usr/bin/shorewalld --listen-prom :9748 --netns %i $SHOREWALLD_ARGS
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/shorewalld /var/lib/shorewalld /var/log
PrivateTmp=yes

# Capabilities:
#   CAP_NET_ADMIN  — nft writes inside netns; also required to bind
#                    the nfnetlink_log socket when LOG_DISPATCH=shorewalld
#   CAP_SYS_ADMIN  — setns(2) into named network namespaces
#   CAP_NET_RAW    — conservative, not strictly required
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN

RuntimeDirectory=shorewalld-%i
RuntimeDirectoryMode=0750
StateDirectory=shorewalld
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
